# etc/supervisor/conf.d/rodan.conf.in
@IF_DIVA@[fcgi-program:iipserver]
@IF_DIVA@command=@IIPSRV_FCGI@
@IF_DIVA@socket=unix:///tmp/%(program_name)s_1.sock
@IF_DIVA@user=@WWW_USER@
@IF_DIVA@autostart=true
@IF_DIVA@autorestart=unexpected
@IF_DIVA@redirect_stderr=true
@IF_DIVA@redirect_stdout=true
@IF_DIVA@# KEEP FILESYSTEM_PREFIX same as the media path
@IF_DIVA@environment=JPEG_QUALITY='100',MAX_LAYERS='4',MAX_CVT='7000',FILESYSTEM_PREFIX='@RODAN_DATA_DIR@'

[program:rodan-celery]
command=@RODAN_APP_DIR@celery_start.sh
environment=PYTHON_EGG_CACHE="/tmp"
directory=@RODAN_APP_DIR@
chown=@WWW_USER@:@WWW_GROUP@
user=@WWW_USER@
autostart=true
autorestart=true
redirect_stderr=true
redirect_stdout=true
# Give it enough time to finish current task.
stopwaitsecs=600

@IF_SERVER@[program:rodan]
@IF_SERVER@command=@RODAN_APP_DIR@gunicorn_start.sh
@IF_SERVER@environment=PYTHON_EGG_CACHE="/tmp",DJANGO_SETTINGS_MODULE="rodan.settings"
@IF_SERVER@directory=@RODAN_APP_DIR@
@IF_SERVER@chown=@WWW_USER@:@WWW_GROUP@
@IF_SERVER@user=@WWW_USER@
@IF_SERVER@autostart=true
@IF_SERVER@autorestart=true
@IF_SERVER@redirect_stderr=true
@IF_SERVER@redirect_stdout=true

@IF_SERVER@[program:rodan-websocket]
@IF_SERVER@command=@RODAN_APP_DIR@websocket_start.sh
@IF_SERVER@environment=PYTHON_EGG_CACHE="/tmp",DJANGO_SETTINGS_MODULE="rodan.settings"
@IF_SERVER@directory=@RODAN_APP_DIR@
@IF_SERVER@chown=@WWW_USER@:@WWW_GROUP@
@IF_SERVER@user=@WWW_USER@
@IF_SERVER@autostart=true
@IF_SERVER@autorestart=true
@IF_SERVER@redirect_stderr=true
@IF_SERVER@redirect_stdout=true
@IF_SERVER@stopwaitsecs=10
