# etc/supervisor/conf.d/rodan-celery.conf.in
[program:rodan-celery]
command=@RODAN_APP_DIR@celery_start.sh
environment=PYTHON_EGG_CACHE="/tmp"
directory=@RODAN_APP_DIR@
chown=www-data:www-data
user=www-data
autostart=true
autorestart=true
redirect_stderr=true
redirect_stdout=true