# etc/supervisor/conf.d/rodan-websocket.conf.in
[program:rodan]
command=@RODAN_APP_DIR@websocket_start.sh
environment=PYTHON_EGG_CACHE="/tmp",DJANGO_SETTINGS_MODULE="rodan.settings"
directory=@RODAN_APP_DIR@
chown=@WWW_USER@:@WWW_GROUP@
user=@WWW_USER@
autostart=true
autorestart=true
redirect_stderr=true
redirect_stdout=true
